# services/awww-rotate-once.service
[Unit]
Description=Set a new wallpaper via awww

[Service]
Type=oneshot
ExecStart=%h/.local/bin/awww-rotate-once.sh
